# apps/indexer/helix-indexer.service
# ===== MAINNET leaderboard indexer (chainid 1776) =====
# Indexes PredictionBook BetOpened/BetSettled -> SQLite -> leaderboard.json, published to the
# helix-leaderboard CF Pages project. SEPARATE from the keeper: its own RPC budget + memory cap so a
# backfill burst can never throttle the liveness-critical settler. EnvironmentFile=.env.mainnet.
[Unit]
Description=helix-predict leaderboard indexer — MAINNET (BetOpened/BetSettled -> leaderboard.json)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helixkeeper
Group=helixkeeper
WorkingDirectory=/opt/helix-predict/apps/indexer
# HOME + npm cache under the one ReadWritePath so `npx wrangler` (the publish step) can cache/run
# under ProtectHome=true. First publish downloads wrangler here; later ones are cached.
Environment=HOME=/opt/helix-predict/apps/indexer
Environment=NPM_CONFIG_CACHE=/opt/helix-predict/apps/indexer/.npm
EnvironmentFile=/opt/helix-predict/apps/indexer/.env.mainnet
ExecStart=/usr/local/bin/node --experimental-sqlite src/indexer.js
Restart=always
RestartSec=10
MemoryMax=256M
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/helix-predict/apps/indexer

[Install]
WantedBy=multi-user.target
